# fluxmonitor/data/rapi/udev-99-flux.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="ffff", RUN+="/usr/bin/killall -s SIGUSR2 'flux: UsbService'"
